# data/io.github.pureblueos.purestore.service.in
[Unit]
Description=Store Service
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=io.github.pureblueos.purestore
ExecStart=@bindir@/purestore --no-window
